# software/ros_ws/src/perseus_lite_screen/scripts/perseus-lite-screen.service
[Unit]
Description=Perseus Lite Qt EGLFS map screen (kiosk on DP-1)
Documentation=https://github.com/ROAR-QUTRC/perseus-v2
After=network-online.target
Wants=network-online.target
# Don't run while a display manager owns DRM master.
Conflicts=gdm.service display-manager.service

[Service]
Type=simple
# Root needed to acquire DRM master without an active seat session.
User=root
# At cold boot under multi-user.target the Tegra DRM (display) and
# Nvidia DRM (GPU) modules are not pulled in by anything — Xorg/GDM
# would normally do that. Force-load them here so /dev/dri/card1
# exists before the binary opens it. `-q` silences nonexistent module
# errors; the leading `-` lets the unit continue regardless.
ExecStartPre=-/usr/sbin/modprobe -q tegra-drm
ExecStartPre=-/usr/sbin/modprobe -q nvidia-drm modeset=1
ExecStart=/usr/local/bin/perseus-lite-screen-kiosk
Restart=on-failure
RestartSec=5
TimeoutStopSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
